# rootfs/etc/services.d/mariadb/run
#!/bin/sh

if [ -f ./setup ];then
  ./setup
fi
if [ "$?" = "0" ];then
  exec 2>&1
  exec mysqld
else
  exit 1
fi